# cron/com.user.municipal-bid-agent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.municipal-bid-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{PROJECT_DIR}}/scripts/run_daily.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{PROJECT_DIR}}</string>

    <!-- 毎日 {{HOUR}}:{{MINUTE}} に起動 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>{{HOUR}}</integer>
        <key>Minute</key>
        <integer>{{MINUTE}}</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>{{PROJECT_DIR}}/logs/launchd.out</string>
    <key>StandardErrorPath</key>
    <string>{{PROJECT_DIR}}/logs/launchd.err</string>

    <!-- 起動時に missed run があっても1回だけ実行 -->
    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
